# resources/ac.conf
[program:ac]
autorestart=true
autostart=true
directory=/opt/ac-remote
command=/opt/ac-remote/ac
user=pi
redirect_stderr=true
stdout_logfile=/dev/null